# conf/docker/etc/supervisord.conf
[supervisord]
user=root
nodaemon=true

[unix_http_server]
file = /tmp/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 5MB
logfile_backups = 1
pidfile = /var/run/supervisord.pid

[supervisorctl]
serverurl = unix:///tmp/supervisord.sock

[program:webxiang]
command=uwsgi --socket 0.0.0.0:8080 --wsgi-file webxiang/wsgi.py -L --log-4xx --log-5xx --log-master --uid www-data --gid users
environment=DJANGO_SETTINGS_MODULE=run.settings_docker
stopsignal=KILL
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
